# crontab
# m h  dom mon dow   command

5,20,35,50 * * * * /home/unity-merger/bin/fire-tarmac --tarmac-path=/home/unity-merger/tarmac-jenkins /home/unity-merger/tarmac-projects
